# ansible/roles/gunicorn/templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn-{{ project_name }}.sock
SocketUser=nginx

[Install]
WantedBy=sockets.target